# systemd/agriapp-capture.service
[Unit]
Description=AgriApp Camera Capture Loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/fra/agri-app
ExecStart=/home/fra/agri-app/scripts/run_capture.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
